# docker/s6-overlay/s6-rc.d/queue/run
#!/command/execlineb -P
/usr/local/bin/php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --timeout=60 --memory=128
